# mon-services.service
[Unit]
Description=mon-services

[Service]
Type=oneshot
ExecStart=/www/mon-services/mon-services.sh

[Install]
WantedBy=multi-user.target